# scripts/systemd/agent-browser-idle-reaper.service
[Unit]
Description=Clean up idle agent-browser daemons

[Service]
Type=oneshot
ExecStart=%h/summitflow/scripts/agent-browser-idle-reaper.js
