# scripts/systemd/codex-session-sync.timer
[Unit]
Description=Run Codex transcript sync every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=15s
Unit=codex-session-sync.service
Persistent=true

[Install]
WantedBy=timers.target
