# scripts/systemd/codex-session-sync.service
[Unit]
Description=Codex Transcript Sync to Agent Hub
After=network.target

[Service]
Type=oneshot
ExecStart=%h/summitflow/scripts/codex-session-sync.py --scan --verbose
StandardOutput=journal
StandardError=journal
