# scripts/systemd/summitflow-host-guardian.service
[Unit]
Description=Independent host, storage, backup, and infrastructure health guard
After=local-fs.target docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/summitflow-host-guardian check
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=15min
StandardOutput=journal
StandardError=journal

