# scripts/systemd/summitflow-celery-beat.service
[Unit]
Description=SummitFlow Celery Beat Scheduler
After=network.target postgresql.service redis-server.service summitflow-celery.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=%h/summitflow/backend
Environment="PATH=%h/summitflow/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=%h"
EnvironmentFile=-%h/summitflow/backend/.env
ExecStart=%h/summitflow/backend/.venv/bin/celery -A app.celery_app beat --loglevel=info
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=summitflow-beat

[Install]
WantedBy=default.target
