# scripts/systemd/summitflow-backend.service
[Unit]
Description=SummitFlow Backend (FastAPI)
After=network.target

[Service]
Type=simple
WorkingDirectory=__SUMMITFLOW_ROOT__/backend
Environment="PATH=%h/.local/bin:__SUMMITFLOW_ROOT__/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=%h"
EnvironmentFile=-%h/.env.local
Environment="PYTHONUNBUFFERED=1"
Environment="HATCHET_CLIENT_HOST_PORT=127.0.0.1:7070"
Environment="HATCHET_CLIENT_TLS_STRATEGY=none"
ExecStart=__SUMMITFLOW_ROOT__/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=10
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=summitflow-backend

[Install]
WantedBy=default.target
