# scripts/systemd/summitflow-host-maintenance.timer
[Unit]
Description=Run conservative host maintenance daily

[Timer]
OnCalendar=*-*-* 05:15:00
RandomizedDelaySec=10min
Persistent=true
Unit=summitflow-host-maintenance.service

[Install]
WantedBy=timers.target

