# scripts/systemd/summitflow-nvme-short-test.timer
[Unit]
Description=Run an NVMe short self-test weekly

[Timer]
OnCalendar=Sat *-*-* 03:00:00
RandomizedDelaySec=10min
Persistent=true
Unit=summitflow-nvme-short-test.service

[Install]
WantedBy=timers.target

