# scripts/systemd/summitflow-terminal.service
[Unit]
Description=SummitFlow Terminal Service
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=%h/summitflow/services/terminal
ExecStart=%h/summitflow/services/terminal/.venv/bin/python -m terminal
EnvironmentFile=-%h/.env.local
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
