# scripts/systemd/summitflow-terminal-frontend.service
[Unit]
Description=SummitFlow Terminal Frontend (Next.js)
After=network.target summitflow-terminal.service
Wants=summitflow-terminal.service

[Service]
Type=simple
WorkingDirectory=%h/terminal/frontend
Environment="NODE_ENV=production"
Environment="PORT=3002"
# Use bash -lc to source .bashrc/.profile for NVM/pnpm path resolution
ExecStart=/bin/bash -lc "pnpm start --hostname 0.0.0.0 --port 3002"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=summitflow-terminal-frontend

[Install]
WantedBy=default.target
