# scripts/systemd/tmux-agent-session-sync.service
[Unit]
Description=Sync external tmux agent sessions into Agent Hub
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env python3 %h/summitflow/scripts/tmux-agent-session-sync.py

